# gu-hub/deployment/deb/prerm

systemctl stop gu-hub
systemctl disable gu-hub


