# gu-provider/deployment/deb/prerm

systemctl stop gu-provider
systemctl disable gu-provider


